# templates/systemd/domain-el.service.erb
[Unit]
Description = GlassFish Server <%= scope.lookupvar('glassfish::version') %> - <%= @domain_name %>
After = syslog.target network-online.target

[Service]
User=<%= @runuser %>
ExecStart = <%= scope.lookupvar('glassfish::glassfish_dir') %>/bin/asadmin start-domain <%= @domain_name %>
ExecStop = <%= scope.lookupvar('glassfish::glassfish_dir') %>/bin/asadmin stop-domain <%= @domain_name %>
ExecReload = <%= scope.lookupvar('glassfish::glassfish_dir') %>/bin/asadmin restart-domain <%= @domain_name %>
Type = forking

[Install]
WantedBy = multi-user.target

