# templates/systemd/instance.service.erb
# This file is managed by Puppet
<% $domain_name = @svc_name.gsub('glassfish_', '') -%>
[Unit]
Description = GlassFish Server <%= scope.lookupvar('glassfish::version') %> - <%= $domain_name %>
After = syslog.target network-online.target

[Service]
User=<%= @user %>
<% if @start_timeout -%>
TimeoutStartSec = <%= @start_timeout %>
<% end -%>
ExecStart = <%= scope.lookupvar('glassfish::glassfish_dir') %>/bin/asadmin start-domain <%= $domain_name %>
ExecStop = <%= scope.lookupvar('glassfish::glassfish_dir') %>/bin/asadmin stop-domain <%= $domain_name %>
ExecReload = <%= scope.lookupvar('glassfish::glassfish_dir') %>/bin/asadmin restart-domain <%= $domain_name %>
Type = forking

[Install]
WantedBy = multi-user.target

