# templates/glassfish-init-cluster-el.erb
#!/bin/bash
#
# chkconfig: 2345 99 35
# description: Starts and stops the Glassfish Server \
#              used to provide <%= @domain %> Domain.
#

# source function library
. /etc/init.d/functions

RETVAL=0
GLASSFISH_BIN="<%= scope.lookupvar('glassfish::glassfish_dir') %>/glassfish/bin"

# Setup Cluster and User values
CLUSTER="<%= @cluster_name %>"
DAS_PORT="<%= @das_port %>"
USER="<%= @runuser %>"

start() {
   echo -n "Starting glassfish cluster $CLUSTER: "
   daemon --user $USER --pidfile /var/run/glassfish_$CLUSTER.pid "$GLASSFISH_BIN/asadmin --port $DAS_PORT start-cluster $CLUSTER"
   RETVAL=$?
   if [ $RETVAL -eq 0 ]; then
        touch /var/lock/subsys/glassfish_$CLUSTER
   fi
   echo
}

stop() {
    echo -n "Stopping glassfish cluster $CLUSTER: "
    su - $USER -c '$GLASSFISH_BIN/asadmin stop-cluster --port $DAS_PORT $CLUSTER'
    RETVAL=$?
    [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/glassfish_$CLUSTER && success || failure
    echo
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
esac

exit $RETVAL
